# scripts/maxids.service
[Unit]
Description=Mini-AD Max UID/GID JSON Service
After=network.target

[Service]
# Run as root (needed for port 80)
User=root
Group=root

# Path to Python and your script
ExecStart=/usr/bin/python3 /opt/maxids_service.py

# Restart on failure
Restart=on-failure
RestartSec=5s

# Working directory
WorkingDirectory=/opt

[Install]
WantedBy=multi-user.target
